# deploy/vocab-trainer-watcher.path
[Path]
Unit=vocab-trainer-watcher.service
PathChanged=/opt/vocab-trainer/vocab-trainer

[Install]
WantedBy=multi-user.target
