# deploy/vocab-trainer.service
# systemd service for vocab-trainer
#
# Install:
#   sudo cp deploy/vocab-trainer.service /etc/systemd/system/vocab-trainer.service
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now vocab-trainer
#
# The service automatically restarts whenever the binary is replaced
# (WatchPaths triggers a restart on inotify events for the binary file).

[Unit]
Description=Vocabulary Trainer
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/vocab-trainer

# DB stored next to the binary; override with DB_PATH env var if preferred
Environment=DB_PATH=/opt/vocab-trainer/data/vocab.db
Environment=PORT=8080

# TTS: remove these lines to disable audio generation
# Requires: cd /opt/vocab-trainer && python3 -m venv tts-venv && tts-venv/bin/pip install -r requirements.txt
Environment=TTS_SCRIPT=/opt/vocab-trainer/generate.py
Environment=VENV_PYTHON=/opt/vocab-trainer/tts-venv/bin/python3
Environment=AUDIO_DIR=/opt/vocab-trainer/data/audio

ExecStart=/opt/vocab-trainer/vocab-trainer
Restart=on-failure
RestartSec=3s

# Restart automatically when the binary is updated via rsync
WatchPaths=/opt/vocab-trainer/vocab-trainer

StandardOutput=journal
StandardError=journal
SyslogIdentifier=vocab-trainer

[Install]
WantedBy=multi-user.target
